# config/autostart/diodon.desktop
[Desktop Entry]
Name=Diodon
Comment=Clipboard Manager
Icon=diodon
Exec=diodon
Terminal=false
Type=Application
