# config/autostart/nm-applet.desktop
[Desktop Entry]
Type=Application
Name=Network Manager
Comment=Network Manager applet
Exec=nm-applet
StartupNotify=false
Terminal=false
Hidden=false
